# src/echodataflow/services/deploy_prefect_worker.service
[Unit]
Description=Deploy Prefect worker in work-pool local
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/exouser
Environment=PATH=/home/exouser/miniforge3/condabin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/exouser/.config/echodataflow/services.env
ExecStart=/bin/sh -lc 'exec "$MAMBA_BIN" run -n "$ECHODATAFLOW_ENV" prefect worker start --pool "$PREFECT_POOL"'
Restart=always
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echodataflow_prefect_worker

[Install]
WantedBy=multi-user.target
