# src/echodataflow/services/deploy_prefect_server.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.echodataflow.prefect-server</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>WorkingDirectory</key>
  <string>/Users/feresa/code_git/echodataflow</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>
ENV_FILE=/Users/feresa/.config/echodataflow/services.env
[[ -f "$ENV_FILE" ]] || { echo "Missing $ENV_FILE"; exit 1; }

set -a
source "$ENV_FILE"
set +a

DB_PATH="${PREFECT_DB_PATH:-$HOME/.prefect/prefect.db}"
AUTO_RESET="${PREFECT_DB_AUTORESET_ON_FAILURE:-false}"

if [[ -f "$DB_PATH" ]]; then
  DB_CHECK="$(sqlite3 "$DB_PATH" "PRAGMA quick_check;" 2>/dev/null || true)"
  if [[ "$DB_CHECK" != "ok" ]]; then
    echo "Prefect DB health check failed for $DB_PATH: ${DB_CHECK:-no output}"
    if [[ "$AUTO_RESET" == "true" ]]; then
      ts="$(date +%Y%m%d_%H%M%S)"
      backup_dir="$HOME/.prefect/db-backups/$ts"
      mkdir -p "$backup_dir"
      mv "${DB_PATH}"* "$backup_dir"/ 2>/dev/null || true
      echo "Archived DB files to $backup_dir"
    else
      echo "Set PREFECT_DB_AUTORESET_ON_FAILURE=true in services.env to auto-archive and reset."
      exit 1
    fi
  fi
fi

exec "$MAMBA_BIN" run -n "$ECHODATAFLOW_ENV" prefect server start --host 127.0.0.1 --port 4200
    </string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/feresa/.local/var/log/echodataflow/prefect-server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/feresa/.local/var/log/echodataflow/prefect-server.err.log</string>
</dict>
</plist>
