# src/echodataflow/services/deploy_prefect_server.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.echodataflow.prefect-server</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/feresa/code_git/echodataflow</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>ENV_FILE=/Users/feresa/.config/echodataflow/services.env; [[ -f "$ENV_FILE" ]] || { echo "Missing $ENV_FILE"; exit 1; }; set -a; source "$ENV_FILE"; set +a; exec "$MAMBA_BIN" run -n "$ECHODATAFLOW_ENV" prefect server start --host 127.0.0.1 --port 4200</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/feresa/.local/var/log/echodataflow/prefect-server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/feresa/.local/var/log/echodataflow/prefect-server.err.log</string>
</dict>
</plist>
